# deploy/maxek-erp.service
[Unit]
Description=MAXEK ERP Web Application (Flask + Gunicorn)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/maxek-erp
EnvironmentFile=-/var/www/maxek-erp/.env
Environment=PATH=/var/www/maxek-erp/venv/bin
Environment=MAXEK_SKIP_DEMO_SEED=1
ExecStart=/var/www/maxek-erp/venv/bin/gunicorn --workers 2 --bind 127.0.0.1:8000 --timeout 120 wsgi:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
